# example/pico433_serial_mqtt.timer
[Timer]
OnBootSec=30sec
[Install]
WantedBy=timers.target

